# packages/sysutils/odroid-remote/system.d/odroid-remote.service
[Unit]
Description=ODROID Remote Control Setup

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/odroid_remote /etc/odroid-remote/odroid_remote.conf

[Install]
WantedBy=multi-user.target
